# PC-APP-V1/pc-rpi-agent/test/service_test
#! /bin/sh
### BEGIN INIT INFO

#output=$((./service_test.py start) 2>&1)
$(./service_test.py start)

echo $output

echo "----- exit value ------"
echo $?